# install/qr-reader-agent/systemd/viaaccess-qr-agent.service
[Unit]
Description=ViaAccess QR Reader appliance agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=viaaccess
Group=viaaccess
SupplementaryGroups=input gpio
EnvironmentFile=-/etc/viaaccess-qr-reader/env
# --hid-auto: USB keyboard-wedge via evdev (stdin alone does not receive HID under systemd)
ExecStart=/usr/local/bin/viaaccess-qr-agent --config /etc/viaaccess-qr-reader/config.json --hid-auto
Restart=always
RestartSec=3
NoNewPrivileges=true
# Any DeviceAllow switches the cgroup to a closed device policy — list HID + GPIO.
DevicePolicy=closed
DeviceAllow=char-input r
DeviceAllow=/dev/input/event* r
DeviceAllow=/dev/gpiochip0 rw
DeviceAllow=char-gpio rw

[Install]
WantedBy=multi-user.target
